# systemd/qopenhd.service
[Unit]
Description=QOpenHD
After=multi-user.target

[Service]
Type=simple
Environment="QT_QPA_EGLFS_KMS_ATOMIC=1"
Environment="QT_QPA_EGLFS_KMS_CONFIG=/usr/local/share/qopenhd/rpi_qt_eglfs_kms_config.json"
ExecStart=/usr/local/bin/QOpenHD
User=root
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

